# low-battery-notify/low-battery-notify.service
[Unit]
Description=Low battery notify

[Service]
Type=oneshot
ExecStart=/home/flygrounder/.local/bin/low-battery-notify.py
